# chopin/supervisor/celery_dev.conf
; =============================================================================
; celery beat supervisor configuration
;
; sudo ln -s `pwd`/celery_dev.conf /etc/supervisor/conf.d/
; =============================================================================

[program:celery-dev]
command = /vol/ocve3/webroot/envs/chopin_dev/celery -A chopin beat --loglevel=INFO
directory = /vol/ocve3/webroot/dev/django/chopin
stdout_logfile = /vol/ocve3/webroot/dev/django/chopin/logs/celery_beat.log
stderr_logfile = /vol/ocve3/webroot/dev/django/chopin/logs/celery_beat.log
autostart = true
autorestart = true
startsecs = 10
